# templates/default/runsvdir-start.service.erb
[Unit]
Description=<%= @project_name %> Runit Process Supervisor
After=network.target auditd.service

[Service]
ExecStart=<%= @install_path %>/embedded/bin/runsvdir-start
Restart=always

[Install]
WantedBy=multi-user.target
